# box/systemd/box-prune.timer
[Unit]
Description=box weekly prune timer

[Timer]
OnCalendar=weekly
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target
