# box/systemd/box-snapshot.timer
[Unit]
Description=box hourly snapshot timer

[Timer]
OnCalendar=hourly
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
