# box/systemd/box-prune.service
[Unit]
Description=box weekly restic prune

[Service]
Type=oneshot
ExecStart=/usr/local/bin/box prune
